# .local/scripts/new_machine_setup/arch_update/arch_auto_update.timer
[Unit]
Description=Timer to trigger automated system update

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target

